# microblog.conf
[program:microblog]
command: /usr/bin/python microblog.py runserver
directory: /home/ubuntu/FlaskMicroblog
autostart: true