# microblog.conf
[program:microblog]
command: /usr/bin/python -m gevent_wrapper
directory: /home/ubuntu/FlaskMicroblog
autostart: true