# microblog.conf
[program:microblog]
command: /usr/bin/python -m microblog
directory: /home/ubuntu/FlaskMicroblog
autostart: true
environment=MICROBLOG_CONFIG="/home/ubuntu/FlaskMicroblog/config.py"